# crontab
* * * * * /Users/oya/git-rep/crontab/autoGitPull.sh
